# firewall.conf
description "run firewall init when eth0 comes up"

start on (local-filesystems
      and remote-filesystems
      and net-device-up IFACE=eth0)

task

exec /usr/local/sbin/iptables.sh
